# docs/systemd/case-weather-cache.service
[Unit]
Description=Case Weather cache sync
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/your-app
EnvironmentFile=/opt/your-app/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash /opt/your-app/scripts/weather_cache_sync.sh

[Install]
WantedBy=multi-user.target
