# docs/systemd/case-weather-cache.timer
[Unit]
Description=Run Case Weather cache sync every 30 minutes
ConditionPathExists=|!/opt/your-app/deployments/activation-in-progress
ConditionPathExists=|/run/case-weather/activation-permit

[Timer]
OnActiveSec=30min
OnUnitActiveSec=30min
AccuracySec=1s
Unit=case-weather-cache.service

[Install]
WantedBy=timers.target
