# docs/systemd/case-weather-cache-bootstrap.service
[Unit]
Description=DEPRECATED - Case Weather cache bootstrap must not be enabled
RefuseManualStart=yes
ConditionPathExists=/run/case-weather/deprecated-bootstrap-never-enable

[Service]
Type=oneshot
User=case-weather
Group=case-weather
UMask=0077
# 兼容旧主机清理时保留一个无网络、无写入的安全占位单元。
ExecStart=/usr/bin/true
