# docs/systemd/case-weather-cache-bootstrap.timer
[Unit]
Description=Run the first Case Weather cache sync after a full 30 minute wait
ConditionPathExists=|!/opt/your-app/deployments/activation-in-progress
ConditionPathExists=|/run/case-weather/activation-permit

[Timer]
OnActiveSec=30min
AccuracySec=1s
RemainAfterElapse=no
Unit=case-weather-cache.service

[Install]
WantedBy=timers.target
